# 69-hdparm.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sda", RUN+="/usr/bin/hdparm -a 256 -c 3 -M 255 -W 1 -B 255 -S 0 /dev/sda"
